# configs/etc/xdg/autostart/dotfiles.desktop
[Desktop Entry]
Type=Application
Exec=dotfiles
Name=dotfiles
Comment=Copy dotfiles from NTFS home directory
